# deploy/gamlachain.service
[Unit]
Description=GamlaChain Blockchain Server
After=network.target

[Service]
Type=simple
User=gamla
Group=gamla
WorkingDirectory=/opt/gamlachain
Environment=HOST=127.0.0.1
Environment=PORT=8000
Environment=MINING_DIFFICULTY=4
Environment=MINING_REWARD=50.0
Environment=DATA_DIR=/opt/gamlachain/data
Environment=CORS_ORIGINS=*
Environment=SECRET_KEY=
ExecStart=/opt/gamlachain/venv/bin/python3 -m gamla_chain
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
